# scripts/systemd/drofbot-traffic-capture.service
[Unit]
Description=Drofbot VPN Traffic Capture
After=network.target wg-quick@wg0.service
Wants=drofbot-traffic-capture.timer

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /opt/drofbot/scripts/traffic-capture.py
WorkingDirectory=/opt/drofbot
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=multi-user.target
